# infra/systemd/uis-docker.service
[Unit]
Description=Urban Immune System — Docker Compose (TimescaleDB + Qdrant + Kafka)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/wlsdud5035/urban-immune-system
ExecStart=/usr/bin/docker compose up -d timescaledb qdrant
ExecStop=/usr/bin/docker compose stop timescaledb qdrant
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
