# infra/systemd/uis-frontend.service
[Unit]
Description=UIS Next.js Frontend
After=network.target uis-backend.service
Wants=uis-backend.service

[Service]
Type=simple
User=wlsdud5035
Group=wlsdud5035
WorkingDirectory=/home/wlsdud5035/urban-immune-system/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
EnvironmentFile=-/home/wlsdud5035/urban-immune-system/frontend/.env.production
ExecStart=/usr/bin/npm run start -- -H 127.0.0.1 -p 3000
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=120
MemoryMax=1G
CPUQuota=80%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=uis-frontend

[Install]
WantedBy=multi-user.target
