# infra/systemd/uis-scheduler.service
[Unit]
Description=Urban Immune System — APScheduler (3-Layer collection + ensemble + nightly RAG)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=wlsdud5035
Group=wlsdud5035
WorkingDirectory=/home/wlsdud5035/urban-immune-system
EnvironmentFile=/home/wlsdud5035/urban-immune-system/.env
ExecStart=/home/wlsdud5035/urban-immune-system/.venv/bin/python -m pipeline.collectors.scheduler
Restart=on-failure
RestartSec=15
StandardOutput=append:/var/log/uis-scheduler.log
StandardError=append:/var/log/uis-scheduler.log

# 보안 하드닝
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/wlsdud5035/urban-immune-system /var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
