# infra/systemd/uis-backend.service.d/restart-limit.conf
[Unit]
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
RestartSec=10
